# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; PHP-FPM
[program:php-fpm]
command=/usr/local/sbin/php-fpm
priority=999
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log
autorestart=true
autostart=true

; Nginx
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=998
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
autorestart=true
autostart=true

; Laravel Queue Worker (optional - uncomment if using queues)
; [program:laravel-worker]
; process_name=%(program_name)s_%(process_num)02d
; command=php /app/artisan queue:work --sleep=3 --tries=3
; priority=997
; autostart=true
; autorestart=true
; numprocs=1
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/worker.log
; stopwait=3600

; Schedule Runner (optional - uncomment if using scheduling)
; [program:laravel-schedule]
; command=php /app/artisan schedule:run
; autostart=true
; autorestart=true
; stdout_logfile=/var/log/supervisor/schedule.log
; stderr_logfile=/var/log/supervisor/schedule-error.log
